# systemd-timers/poliloom-pipeline.service
[Unit]
Description=PoliLoom download, import and garbage collection pipeline
After=postgresql.service
Requires=postgresql.service

[Service]
Type=oneshot
WorkingDirectory=/root/poliloom
EnvironmentFile=/root/poliloom/.env
ExecStart=docker compose run --rm api poliloom dump-download --output ${WIKIDATA_DUMP_COMPRESSED}
ExecStart=docker compose run --rm api poliloom dump-extract --input ${WIKIDATA_DUMP_COMPRESSED} --output ${WIKIDATA_DUMP_EXTRACTED}
ExecStart=docker compose run --rm api poliloom import-hierarchy --file ${WIKIDATA_DUMP_EXTRACTED}
ExecStart=docker compose run --rm api poliloom import-entities --file ${WIKIDATA_DUMP_EXTRACTED}
ExecStart=docker compose run --rm api poliloom import-politicians --file ${WIKIDATA_DUMP_EXTRACTED}
ExecStart=docker compose run --rm api poliloom embed-entities
ExecStart=docker compose run --rm api poliloom garbage-collect
StandardOutput=journal
StandardError=journal
TimeoutStartSec=259200
