# systemd-timers/poliloom-pipeline.timer
[Unit]
Description=Run PoliLoom pipeline daily at 3 AM
Requires=poliloom-pipeline.service

[Timer]
OnCalendar=daily
Persistent=true
RandomizedDelaySec=3600

[Install]
WantedBy=timers.target