# systemd-timers/poliloom-enrich-stateless.timer
[Unit]
Description=Run PoliLoom stateless enrichment daily

[Timer]
OnCalendar=*-*-* 06:00:00
Persistent=true
RandomizedDelaySec=1800

[Install]
WantedBy=timers.target
